# backend/python-worker/deployment/spot-interruption-handler.service
[Unit]
Description=Spot Instance Interruption Handler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ec2-user/python-worker

# Command
ExecStart=/usr/bin/python3.11 /home/ec2-user/python-worker/spot_interruption_handler.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spot-handler

[Install]
WantedBy=multi-user.target
